# meta-supra/recipes-connectivity/wifi-config/files/wifi-static.service
[Unit]
Description=WiFi static IP configuration
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wifi-static.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
